# files/etc/services.d/mariadb/run
#!/usr/bin/with-contenv bash

exec mysqld_safe --pid-file=$PID_FILE
